# spool/cron/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.puW7V0/crontab installed on Fri Feb  4 10:06:40 2011)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
PATH=/usr/sbin:/usr/bin:/sbin:/bin:/home/chris/Scripts
# m h  dom mon dow   command
0 06 * * * /home/chris/Scripts/sys_backup.sh
