# lib/dpkg/info/apt.conffiles
/etc/apt/apt.conf.d/01autoremove
/etc/apt/apt.conf.d/01ubuntu
/etc/logrotate.d/apt
/etc/cron.daily/apt
