# lib/dpkg/info/aptitude.conffiles
/etc/apt/apt.conf.d/05aptitude
/etc/cron.daily/aptitude
/etc/logrotate.d/aptitude
