# lib/dpkg/info/bind9.list
/.
/etc
/etc/ufw
/etc/ufw/applications.d
/etc/ufw/applications.d/bind9
/etc/apparmor.d
/etc/apparmor.d/force-complain
/etc/apparmor.d/usr.sbin.named
/etc/bind
/etc/bind/bind.keys
/etc/bind/db.0
/etc/bind/db.255
/etc/bind/db.empty
/etc/bind/zones.rfc1918
/etc/bind/db.127
/etc/bind/db.local
/etc/bind/db.root
/etc/bind/named.conf
/etc/bind/named.conf.local
/etc/bind/named.conf.default-zones
/etc/bind/named.conf.options
/etc/ppp
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/bind9
/etc/ppp/ip-down.d
/etc/ppp/ip-down.d/bind9
/etc/network
/etc/network/if-up.d
/etc/network/if-up.d/bind9
/etc/network/if-down.d
/etc/network/if-down.d/bind9
/etc/init.d
/etc/init.d/bind9
/usr
/usr/bin
/usr/sbin
/usr/sbin/named
/usr/sbin/dnssec-keyfromlabel
/usr/sbin/dnssec-dsfromkey
/usr/sbin/dnssec-revoke
/usr/sbin/dnssec-settime
/usr/sbin/arpaname
/usr/sbin/named-journalprint
/usr/sbin/nsec3hash
/usr/sbin/genrandom
/usr/sbin/isc-hmac-fixup
/usr/sbin/ddns-confgen
/usr/share
/usr/share/bind9
/usr/share/bind9/bind9-default.md5sum
/usr/share/man
/usr/share/man/man5
/usr/share/man/man5/named.conf.5.gz
/usr/share/man/man5/rndc.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/named.8.gz
/usr/share/man/man8/dnssec-keyfromlabel.8.gz
/usr/share/man/man8/dnssec-revoke.8.gz
/usr/share/man/man8/dnssec-settime.8.gz
/usr/share/man/man8/isc-hmac-fixup.8.gz
/usr/share/man/man8/named-journalprint.8.gz
/usr/share/man/man8/nsec3hash.8.gz
/usr/share/man/man8/genrandom.8.gz
/usr/share/man/man8/ddns-confgen.8.gz
/usr/share/man/man8/dnssec-dsfromkey.8.gz
/usr/share/man/man1
/usr/share/man/man1/arpaname.1.gz
/usr/share/man/man1/isc-config.sh.1.gz
/usr/share/doc
/usr/share/doc/bind9
/usr/share/doc/bind9/README.Debian.gz
/usr/share/doc/bind9/changelog.Debian.gz
/usr/share/doc/bind9/copyright
/usr/share/doc/bind9/NEWS.Debian.gz
/usr/share/doc/bind9/FAQ.gz
/usr/share/doc/bind9/README.gz
/usr/lib
/var
/var/cache
/var/cache/bind
/var/run
/var/run/named
