# lib/dpkg/info/dbus.list
/.
/var
/var/lib
/var/lib/dbus
/usr
/usr/share
/usr/share/doc
/usr/share/doc/dbus
/usr/share/doc/dbus/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/dbus-send.1.gz
/usr/share/man/man1/dbus-uuidgen.1.gz
/usr/share/man/man1/dbus-monitor.1.gz
/usr/share/man/man1/dbus-daemon.1.gz
/usr/share/man/man1/dbus-cleanup-sockets.1.gz
/usr/share/dbus-1
/usr/share/dbus-1/services
/usr/share/dbus-1/system-services
/usr/bin
/usr/bin/dbus-send
/usr/bin/dbus-monitor
/usr/lib
/usr/lib/dbus-1.0
/etc
/etc/init
/etc/init/dbus.conf
/etc/init.d
/etc/default
/etc/default/dbus
/etc/dbus-1
/etc/dbus-1/session.conf
/etc/dbus-1/system.conf
/etc/dbus-1/system.d
/etc/dbus-1/session.d
/bin
/bin/dbus-daemon
/bin/dbus-cleanup-sockets
/bin/dbus-uuidgen
/lib
/lib/dbus-1.0
/lib/dbus-1.0/dbus-daemon-launch-helper
/usr/share/doc/dbus/AUTHORS
/usr/share/doc/dbus/changelog.gz
/usr/share/doc/dbus/README.gz
/usr/share/doc/dbus/changelog.Debian.gz
/usr/lib/dbus-1.0/dbus-daemon-launch-helper
/etc/init.d/dbus
