# lib/dpkg/info/dovecot-common.list
/.
/etc
/etc/ssl
/etc/ssl/certs
/etc/ssl/private
/etc/ufw
/etc/ufw/applications.d
/etc/ufw/applications.d/dovecot-common
/etc/dovecot
/etc/dovecot/conf.d
/etc/dovecot/auth.d
/etc/dovecot/dovecot-db-example.conf
/etc/dovecot/dovecot-dict-sql-example.conf
/etc/pam.d
/etc/pam.d/dovecot
/etc/default
/etc/default/dovecot
/etc/init.d
/etc/init.d/dovecot
/usr
/usr/sbin
/usr/sbin/dovecot
/usr/sbin/dovecotpw
/usr/share
/usr/share/dovecot
/usr/share/dovecot/dovecot.conf
/usr/share/dovecot/dovecot-ldap.conf
/usr/share/dovecot/dovecot-sql.conf
/usr/share/doc
/usr/share/doc/dovecot-common
/usr/share/doc/dovecot-common/wiki
/usr/share/doc/dovecot-common/wiki/AuthDatabase.CheckPassword.txt.gz
/usr/share/doc/dovecot-common/wiki/AuthDatabase.SQL.txt.gz
/usr/share/doc/dovecot-common/wiki/AuthDatabase.LDAP.AuthBinds.txt
/usr/share/doc/dovecot-common/wiki/AuthDatabase.LDAP.PasswordLookups.txt
/usr/share/doc/dovecot-common/wiki/AuthDatabase.LDAP.Userdb.txt
/usr/share/doc/dovecot-common/wiki/AuthDatabase.LDAP.txt
/usr/share/doc/dovecot-common/wiki/AuthDatabase.Passwd.txt
/usr/share/doc/dovecot-common/wiki/AuthDatabase.PasswdFile.txt
/usr/share/doc/dovecot-common/wiki/IndexFiles.txt.gz
/usr/share/doc/dovecot-common/wiki/AuthDatabase.VPopMail.txt
/usr/share/doc/dovecot-common/wiki/AuthDatabase.txt
/usr/share/doc/dovecot-common/wiki/Authentication.MasterUsers.txt.gz
/usr/share/doc/dovecot-common/wiki/Authentication.PasswordSchemes.txt.gz
/usr/share/doc/dovecot-common/wiki/LDA.txt.gz
/usr/share/doc/dovecot-common/wiki/Authentication.Mechanisms.DigestMD5.txt
/usr/share/doc/dovecot-common/wiki/Authentication.Mechanisms.NTLM.txt
/usr/share/doc/dovecot-common/wiki/Authentication.Mechanisms.Winbind.txt
/usr/share/doc/dovecot-common/wiki/Authentication.Mechanisms.txt
/usr/share/doc/dovecot-common/wiki/Authentication.MultipleDatabases.txt
/usr/share/doc/dovecot-common/wiki/BasicConfiguration.txt.gz
/usr/share/doc/dovecot-common/wiki/Authentication.RestrictAccess.txt
/usr/share/doc/dovecot-common/wiki/Authentication.txt
/usr/share/doc/dovecot-common/wiki/CompilingSource.txt.gz
/usr/share/doc/dovecot-common/wiki/Chrooting.txt
/usr/share/doc/dovecot-common/wiki/Clients.NegativeUIDs.txt
/usr/share/doc/dovecot-common/wiki/Clients.txt
/usr/share/doc/dovecot-common/wiki/CommandLine.txt
/usr/share/doc/dovecot-common/wiki/LDA.Postfix.txt.gz
/usr/share/doc/dovecot-common/wiki/Debugging.Authentication.txt
/usr/share/doc/dovecot-common/wiki/Debugging.ProcessTracing.txt
/usr/share/doc/dovecot-common/wiki/Debugging.Rawlog.txt
/usr/share/doc/dovecot-common/wiki/Debugging.Thunderbird.txt
/usr/share/doc/dovecot-common/wiki/Design.Indexes.MainIndex.txt.gz
/usr/share/doc/dovecot-common/wiki/LDA.Sieve.txt.gz
/usr/share/doc/dovecot-common/wiki/Design.Indexes.MailIndexApi.txt
/usr/share/doc/dovecot-common/wiki/Design.Processes.txt.gz
/usr/share/doc/dovecot-common/wiki/HowTo.PopBSMTPAndDovecot.txt.gz
/usr/share/doc/dovecot-common/wiki/Design.Indexes.txt
/usr/share/doc/dovecot-common/wiki/Design.MailProcess.txt
/usr/share/doc/dovecot-common/wiki/LDA.Sieve.CMU.txt.gz
/usr/share/doc/dovecot-common/wiki/Design.txt
/usr/share/doc/dovecot-common/wiki/DomainLost.txt
/usr/share/doc/dovecot-common/wiki/FindMailLocation.txt
/usr/share/doc/dovecot-common/wiki/FinishBasicConfiguration.txt
/usr/share/doc/dovecot-common/wiki/HowTo.EximAndDovecotSASL.txt
/usr/share/doc/dovecot-common/wiki/HowTo.PopRelay.txt.gz
/usr/share/doc/dovecot-common/wiki/LDA.Sieve.Dovecot.txt.gz
/usr/share/doc/dovecot-common/wiki/HowTo.PostfixAndDovecotSASL.txt
/usr/share/doc/dovecot-common/wiki/HowTo.Rootless.txt
/usr/share/doc/dovecot-common/wiki/Logging.txt.gz
/usr/share/doc/dovecot-common/wiki/HowTo.txt
/usr/share/doc/dovecot-common/wiki/Migration.txt.gz
/usr/share/doc/dovecot-common/wiki/InetdInstall.txt
/usr/share/doc/dovecot-common/wiki/Iptables.txt
/usr/share/doc/dovecot-common/wiki/LDA.Exim.txt
/usr/share/doc/dovecot-common/wiki/LDA.Indexing.txt
/usr/share/doc/dovecot-common/wiki/NFS.txt.gz
/usr/share/doc/dovecot-common/wiki/LDA.Qmail.txt
/usr/share/doc/dovecot-common/wiki/LDA.Sendmail.txt
/usr/share/doc/dovecot-common/wiki/LoginProcess.txt.gz
/usr/share/doc/dovecot-common/wiki/MailLocation.Maildir.txt.gz
/usr/share/doc/dovecot-common/wiki/MailLocation.Mbox.txt.gz
/usr/share/doc/dovecot-common/wiki/MailLocation.txt.gz
/usr/share/doc/dovecot-common/wiki/MDA.txt
/usr/share/doc/dovecot-common/wiki/MTA.txt
/usr/share/doc/dovecot-common/wiki/MailLocation.LocalDisk.txt
/usr/share/doc/dovecot-common/wiki/MailboxFormat.Maildir.txt.gz
/usr/share/doc/dovecot-common/wiki/MailboxFormat.txt.gz
/usr/share/doc/dovecot-common/wiki/MailLocation.SharedDisk.txt
/usr/share/doc/dovecot-common/wiki/Quota.1.1.txt.gz
/usr/share/doc/dovecot-common/wiki/MailboxFormat.Cydir.txt
/usr/share/doc/dovecot-common/wiki/MailboxFormat.MH.txt
/usr/share/doc/dovecot-common/wiki/MailboxFormat.dbox.txt.gz
/usr/share/doc/dovecot-common/wiki/MboxProblems.txt.gz
/usr/share/doc/dovecot-common/wiki/MailboxFormat.mailstore.txt
/usr/share/doc/dovecot-common/wiki/Migration.Courier.txt.gz
/usr/share/doc/dovecot-common/wiki/MailboxFormat.mbx.txt
/usr/share/doc/dovecot-common/wiki/Migration.UW.txt.gz
/usr/share/doc/dovecot-common/wiki/ManageSieve.Clients.txt
/usr/share/doc/dovecot-common/wiki/ManageSieve.Install.txt.gz
/usr/share/doc/dovecot-common/wiki/ManageSieve.Install.1.0.txt
/usr/share/doc/dovecot-common/wiki/Migration.BincIMAP.txt.gz
/usr/share/doc/dovecot-common/wiki/Migration.MailFormat.txt.gz
/usr/share/doc/dovecot-common/wiki/ManageSieve.txt
/usr/share/doc/dovecot-common/wiki/MboxLocking.txt
/usr/share/doc/dovecot-common/wiki/Namespaces.txt.gz
/usr/share/doc/dovecot-common/wiki/POP3Server.txt.gz
/usr/share/doc/dovecot-common/wiki/Plugins.Expire.txt.gz
/usr/share/doc/dovecot-common/wiki/Migration.Cyrus.txt
/usr/share/doc/dovecot-common/wiki/Migration.Linuxconf.txt
/usr/share/doc/dovecot-common/wiki/Plugins.FTS.Squat.txt.gz
/usr/share/doc/dovecot-common/wiki/Migration.Teapop.txt
/usr/share/doc/dovecot-common/wiki/Quota.FS.txt.gz
/usr/share/doc/dovecot-common/wiki/Migration.Vm-pop3d.txt
/usr/share/doc/dovecot-common/wiki/UserIds.txt.gz
/usr/share/doc/dovecot-common/wiki/MissingMailboxes.txt
/usr/share/doc/dovecot-common/wiki/PasswordDatabase.PAM.txt.gz
/usr/share/doc/dovecot-common/wiki/OSCompatibility.txt
/usr/share/doc/dovecot-common/wiki/Variables.txt.gz
/usr/share/doc/dovecot-common/wiki/PasswordDatabase.BSDAuth.txt
/usr/share/doc/dovecot-common/wiki/PasswordDatabase.ExtraFields.AllowNets.txt
/usr/share/doc/dovecot-common/wiki/PasswordDatabase.ExtraFields.Host.txt
/usr/share/doc/dovecot-common/wiki/PasswordDatabase.ExtraFields.NoDelay.txt
/usr/share/doc/dovecot-common/wiki/PasswordDatabase.ExtraFields.NoLogin.txt
/usr/share/doc/dovecot-common/wiki/Plugins.Virtual.txt.gz
/usr/share/doc/dovecot-common/wiki/PasswordDatabase.ExtraFields.User.txt
/usr/share/doc/dovecot-common/wiki/PasswordDatabase.ExtraFields.txt
/usr/share/doc/dovecot-common/wiki/Quota.Dict.txt.gz
/usr/share/doc/dovecot-common/wiki/PasswordDatabase.Shadow.txt
/usr/share/doc/dovecot-common/wiki/PasswordDatabase.txt
/usr/share/doc/dovecot-common/wiki/PerformanceTuning.txt
/usr/share/doc/dovecot-common/wiki/Plugins.Autocreate.txt
/usr/share/doc/dovecot-common/wiki/Plugins.Convert.txt
/usr/share/doc/dovecot-common/wiki/SystemUsers.txt.gz
/usr/share/doc/dovecot-common/wiki/Plugins.FTS.Lucene.txt
/usr/share/doc/dovecot-common/wiki/Plugins.FTS.Solr.txt
/usr/share/doc/dovecot-common/wiki/RunningDovecot.txt.gz
/usr/share/doc/dovecot-common/wiki/Plugins.FTS.txt
/usr/share/doc/dovecot-common/wiki/Plugins.Lazyexpunge.txt
/usr/share/doc/dovecot-common/wiki/Plugins.Listescape.txt
/usr/share/doc/dovecot-common/wiki/Plugins.MailLog.txt
/usr/share/doc/dovecot-common/wiki/Plugins.MboxSnarf.txt
/usr/share/doc/dovecot-common/wiki/Plugins.Trash.txt
/usr/share/doc/dovecot-common/wiki/Upgrading.1.0.txt.gz
/usr/share/doc/dovecot-common/wiki/Plugins.Zlib.txt
/usr/share/doc/dovecot-common/wiki/Plugins.txt
/usr/share/doc/dovecot-common/wiki/SSL.DovecotConfiguration.txt.gz
/usr/share/doc/dovecot-common/wiki/SharedMailboxes.Permissions.txt.gz
/usr/share/doc/dovecot-common/wiki/Upgrading.1.1.txt.gz
/usr/share/doc/dovecot-common/wiki/Quota.Dirsize.txt
/usr/share/doc/dovecot-common/wiki/Quota.Maildir.txt
/usr/share/doc/dovecot-common/wiki/Quota.txt
/usr/share/doc/dovecot-common/wiki/VirtualUsers.txt.gz
/usr/share/doc/dovecot-common/wiki/SSL.CertificateClientImporting.txt
/usr/share/doc/dovecot-common/wiki/SSL.CertificateCreation.txt
/usr/share/doc/dovecot-common/wiki/TestInstallation.txt.gz
/usr/share/doc/dovecot-common/wiki/SSL.txt
/usr/share/doc/dovecot-common/wiki/Sasl.txt
/usr/share/doc/dovecot-common/wiki/SecurityTuning.txt
/usr/share/doc/dovecot-common/wiki/SharedMailboxes.Public.txt.gz
/usr/share/doc/dovecot-common/wiki/TestPop3Installation.txt.gz
/usr/share/doc/dovecot-common/wiki/SharedMailboxes.Symlinks.txt
/usr/share/doc/dovecot-common/wiki/SharedMailboxes.txt
/usr/share/doc/dovecot-common/wiki/UserDatabase.ExtraFields.txt.gz
/usr/share/doc/dovecot-common/wiki/TimeMovedBackwards.txt
/usr/share/doc/dovecot-common/wiki/Upgrading.1.2.txt
/usr/share/doc/dovecot-common/wiki/Upgrading.txt
/usr/share/doc/dovecot-common/wiki/UserDatabase.NSS.txt
/usr/share/doc/dovecot-common/wiki/UserDatabase.Prefetch.txt
/usr/share/doc/dovecot-common/wiki/UserDatabase.Static.txt
/usr/share/doc/dovecot-common/wiki/UserDatabase.txt
/usr/share/doc/dovecot-common/wiki/WhyDoesItNotWork.txt
/usr/share/doc/dovecot-common/wiki/maildrop.txt
/usr/share/doc/dovecot-common/wiki/mutt.txt
/usr/share/doc/dovecot-common/wiki/uw2dovecot.sh.txt
/usr/share/doc/dovecot-common/wiki/ACL.txt.gz
/usr/share/doc/dovecot-common/wiki/AixPluginsSupport.txt.gz
/usr/share/doc/dovecot-common/wiki/Authentication.Caching.txt.gz
/usr/share/doc/dovecot-common/wiki/Authentication.Kerberos.txt.gz
/usr/share/doc/dovecot-common/wiki/Design.AuthProcess.txt.gz
/usr/share/doc/dovecot-common/wiki/Design.Indexes.Cache.txt.gz
/usr/share/doc/dovecot-common/wiki/Design.Indexes.TransactionLog.txt.gz
/usr/share/doc/dovecot-common/wiki/HowTo.SimpleVirtualInstall.txt.gz
/usr/share/doc/dovecot-common/wiki/MailboxFormat.mbox.txt.gz
/usr/share/doc/dovecot-common/wiki/ManageSieve.Configuration.txt.gz
/usr/share/doc/dovecot-common/wiki/ManageSieve.Troubleshooting.txt.gz
/usr/share/doc/dovecot-common/wiki/PasswordDatabase.ExtraFields.Proxy.txt.gz
/usr/share/doc/dovecot-common/wiki/PostLoginScripting.txt.gz
/usr/share/doc/dovecot-common/wiki/QuickConfiguration.txt.gz
/usr/share/doc/dovecot-common/wiki/SharedMailboxes.Shared.txt.gz
/usr/share/doc/dovecot-common/documentation.txt.gz
/usr/share/doc/dovecot-common/securecoding.txt.gz
/usr/share/doc/dovecot-common/NEWS.gz
/usr/share/doc/dovecot-common/TODO.gz
/usr/share/doc/dovecot-common/sieve.README.gz
/usr/share/doc/dovecot-common/README
/usr/share/doc/dovecot-common/sieve.NEWS.gz
/usr/share/doc/dovecot-common/sieve.ChangeLog.gz
/usr/share/doc/dovecot-common/managesieve.NEWS.gz
/usr/share/doc/dovecot-common/managesieve.ChangeLog.gz
/usr/share/doc/dovecot-common/changelog.Debian.gz
/usr/share/doc/dovecot-common/changelog.gz
/usr/share/doc/dovecot-common/README.Debian
/usr/share/doc/dovecot-common/copyright
/usr/share/doc/dovecot-common/auth-protocol.txt.gz
/usr/share/doc/dovecot-common/thread-refs.txt.gz
/usr/share/doc/dovecot-common/managesieve.README.gz
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/sieve-test.1.gz
/usr/share/man/man1/dovecot.1.gz
/usr/share/man/man1/dovecotpw.1.gz
/usr/share/man/man1/sievec.1.gz
/usr/share/man/man1/sieved.1.gz
/usr/share/man/man1/maildirmake.dovecot.1.gz
/usr/share/man/man8
/usr/share/man/man8/dovecot.8.gz
/usr/lib
/usr/lib/dovecot
/usr/lib/dovecot/dovecot-auth
/usr/lib/dovecot/checkpassword-reply
/usr/lib/dovecot/modules
/usr/lib/dovecot/modules/auth
/usr/lib/dovecot/modules/lda
/usr/lib/dovecot/modules/lda/lib90_sieve_plugin.so
/usr/lib/dovecot/modules/lda/lib90_sieve_plugin.la
/usr/lib/dovecot/modules/lib01_acl_plugin.so
/usr/lib/dovecot/modules/lib01_acl_plugin.la
/usr/lib/dovecot/modules/lib20_autocreate_plugin.so
/usr/lib/dovecot/modules/lib20_autocreate_plugin.la
/usr/lib/dovecot/modules/lib20_convert_plugin.so
/usr/lib/dovecot/modules/lib20_convert_plugin.la
/usr/lib/dovecot/modules/lib20_expire_plugin.so
/usr/lib/dovecot/modules/lib20_expire_plugin.la
/usr/lib/dovecot/modules/lib20_fts_plugin.so
/usr/lib/dovecot/modules/lib20_fts_plugin.la
/usr/lib/dovecot/modules/lib21_fts_squat_plugin.so
/usr/lib/dovecot/modules/lib21_fts_squat_plugin.la
/usr/lib/dovecot/modules/lib02_lazy_expunge_plugin.so
/usr/lib/dovecot/modules/lib02_lazy_expunge_plugin.la
/usr/lib/dovecot/modules/lib20_listescape_plugin.so
/usr/lib/dovecot/modules/lib20_listescape_plugin.la
/usr/lib/dovecot/modules/lib20_mail_log_plugin.so
/usr/lib/dovecot/modules/lib20_mail_log_plugin.la
/usr/lib/dovecot/modules/lib20_mbox_snarf_plugin.so
/usr/lib/dovecot/modules/lib20_mbox_snarf_plugin.la
/usr/lib/dovecot/modules/lib10_quota_plugin.so
/usr/lib/dovecot/modules/lib10_quota_plugin.la
/usr/lib/dovecot/modules/lib11_trash_plugin.so
/usr/lib/dovecot/modules/lib11_trash_plugin.la
/usr/lib/dovecot/modules/lib20_virtual_plugin.so
/usr/lib/dovecot/modules/lib20_virtual_plugin.la
/usr/lib/dovecot/modules/lib20_zlib_plugin.so
/usr/lib/dovecot/modules/lib20_zlib_plugin.la
/usr/lib/dovecot/modules/drac.so
/usr/lib/dovecot/dict
/usr/lib/dovecot/ssl-build-param
/usr/lib/dovecot/deliver
/usr/lib/dovecot/authtest
/usr/lib/dovecot/rawlog
/usr/lib/dovecot/gdbhelper
/usr/lib/dovecot/idxview
/usr/lib/dovecot/listview
/usr/lib/dovecot/logview
/usr/lib/dovecot/maildirlock
/usr/lib/dovecot/threadview
/usr/lib/dovecot/convert-tool
/usr/lib/dovecot/expire-tool
/usr/lib/dovecot/managesieve
/usr/lib/dovecot/managesieve-login
/usr/bin
/usr/bin/sievec
/usr/bin/sieved
/usr/bin/sieve-test
/usr/bin/maildirmake.dovecot
/usr/lib/dovecot/modules/lda/lib01_acl_plugin.so
/usr/lib/dovecot/modules/lda/lib20_autocreate_plugin.so
/usr/lib/dovecot/modules/lda/lib20_convert_plugin.so
/usr/lib/dovecot/modules/lda/lib20_expire_plugin.so
/usr/lib/dovecot/modules/lda/lib20_fts_plugin.so
/usr/lib/dovecot/modules/lda/lib21_fts_squat_plugin.so
/usr/lib/dovecot/modules/lda/lib20_listescape_plugin.so
/usr/lib/dovecot/modules/lda/lib20_mail_log_plugin.so
/usr/lib/dovecot/modules/lda/lib10_quota_plugin.so
/usr/lib/dovecot/modules/lda/lib11_trash_plugin.so
/usr/lib/dovecot/modules/lda/lib20_virtual_plugin.so
