# lib/dpkg/info/fancontrol.list
/.
/usr
/usr/sbin
/usr/sbin/fancontrol
/usr/sbin/pwmconfig
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/fancontrol.8.gz
/usr/share/man/man8/pwmconfig.8.gz
/usr/share/doc
/usr/share/doc/fancontrol
/usr/share/doc/fancontrol/changelog.gz
/usr/share/doc/fancontrol/copyright
/usr/share/doc/fancontrol/changelog.Debian.gz
/usr/share/doc/fancontrol/fancontrol.txt.gz
/etc
/etc/init.d
/etc/init.d/fancontrol
