# lib/dpkg/info/grub-common.list
/.
/usr
/usr/sbin
/usr/sbin/grub-mkdevicemap
/usr/sbin/grub-probe
/usr/sbin/grub-mkconfig
/usr/lib
/usr/lib/grub
/usr/lib/grub/grub-mkconfig_lib
/usr/lib/grub/update-grub_lib
/usr/share
/usr/share/grub
/usr/share/grub/unicode.pf2
/usr/share/grub/ascii.pf2
/usr/share/grub/ascii.h
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/grub-mkisofs.1.gz
/usr/share/man/man1/grub-mkfont.1.gz
/usr/share/man/man1/grub-mkrelpath.1.gz
/usr/share/man/man1/grub-bin2h.1.gz
/usr/share/man/man1/grub-script-check.1.gz
/usr/share/man/man1/grub-editenv.1.gz
/usr/share/man/man1/grub-mkrescue.1.gz
/usr/share/man/man1/grub-mkelfimage.1.gz
/usr/share/man/man1/grub-mkpasswd-pbkdf2.1.gz
/usr/share/man/man8
/usr/share/man/man8/grub-probe.8.gz
/usr/share/man/man8/grub-mkconfig.8.gz
/usr/share/man/man8/grub-mkdevicemap.8.gz
/usr/share/info
/usr/share/info/grub.info.gz
/usr/share/doc
/usr/share/doc/grub-common
/usr/share/doc/grub-common/README
/usr/share/doc/grub-common/TODO
/usr/share/doc/grub-common/AUTHORS
/usr/share/doc/grub-common/THANKS
/usr/share/doc/grub-common/copyright
/usr/share/doc/grub-common/examples
/usr/share/doc/grub-common/examples/grub.cfg
/usr/share/doc/grub-common/NEWS.Debian.gz
/usr/share/doc/grub-common/changelog.gz
/usr/share/doc/grub-common/NEWS.gz
/usr/share/doc/grub-common/changelog.Debian.gz
/usr/share/bug
/usr/share/bug/grub-common
/usr/share/bug/grub-common/presubj
/usr/bin
/usr/bin/grub-mkelfimage
/usr/bin/grub-mkisofs
/usr/bin/grub-mkfont
/usr/bin/grub-mkrelpath
/usr/bin/grub-bin2h
/usr/bin/grub-script-check
/usr/bin/grub-editenv
/usr/bin/grub-mkpasswd-pbkdf2
/usr/bin/grub-mkrescue
/etc
/etc/pm
/etc/pm/sleep.d
/etc/pm/sleep.d/10_grub-common
/etc/grub.d
/etc/grub.d/00_header
/etc/grub.d/30_os-prober
/etc/grub.d/40_custom
/etc/grub.d/10_linux
/etc/grub.d/README
/etc/init.d
/etc/init.d/grub-common
