# lib/dpkg/info/lm-sensors.list
/.
/usr
/usr/sbin
/usr/sbin/isadump
/usr/sbin/isaset
/usr/sbin/sensors-detect
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/isaset.8.gz
/usr/share/man/man8/isadump.8.gz
/usr/share/man/man8/sensors-detect.8.gz
/usr/share/man/man1
/usr/share/man/man1/sensors.1.gz
/usr/share/man/man1/sensors-conf-convert.1.gz
/usr/share/man/man5
/usr/share/man/man5/sensors.conf.5.gz
/usr/share/doc
/usr/share/doc/lm-sensors
/usr/share/doc/lm-sensors/changelog.gz
/usr/share/doc/lm-sensors/donations
/usr/share/doc/lm-sensors/vid
/usr/share/doc/lm-sensors/chips
/usr/share/doc/lm-sensors/chips/MODPARMS.gz
/usr/share/doc/lm-sensors/chips/SUMMARY.gz
/usr/share/doc/lm-sensors/fan-divisors
/usr/share/doc/lm-sensors/applications
/usr/share/doc/lm-sensors/temperature-sensors
/usr/share/doc/lm-sensors/copyright
/usr/share/doc/lm-sensors/examples
/usr/share/doc/lm-sensors/examples/daemon
/usr/share/doc/lm-sensors/examples/daemon/healthd.sh
/usr/share/doc/lm-sensors/examples/hotplug
/usr/share/doc/lm-sensors/examples/hotplug/README.p4b
/usr/share/doc/lm-sensors/examples/hotplug/unhide_ICH_SMBus
/usr/share/doc/lm-sensors/examples/tellerstats
/usr/share/doc/lm-sensors/examples/tellerstats/tellerstats.conf
/usr/share/doc/lm-sensors/examples/tellerstats/tellerstats.sh
/usr/share/doc/lm-sensors/examples/tellerstats/gather.sh
/usr/share/doc/lm-sensors/examples/tellerstats/gnuplotscript.tmpl
/usr/share/doc/lm-sensors/examples/tellerstats/README
/usr/share/doc/lm-sensors/examples/tellerstats/index.shtml
/usr/share/doc/lm-sensors/changelog.Debian.gz
/usr/bin
/usr/bin/sensors
/usr/bin/sensors-conf-convert
/etc
/etc/init.d
/etc/init.d/lm-sensors
