# lib/dpkg/info/php5-common.conffiles
/etc/php5/conf.d/pdo.ini
/etc/cron.d/php5
