# lib/dpkg/info/winbind.conffiles
/etc/logrotate.d/winbind
/etc/default/winbind
/etc/init.d/winbind
