# lib/dpkg/info/winbind.list
/.
/usr
/usr/lib
/usr/lib/samba
/usr/lib/samba/nss_info
/usr/lib/samba/idmap
/usr/lib/samba/idmap/rid.so
/usr/lib/samba/idmap/ad.so
/usr/lib/samba/idmap/adex.so
/usr/lib/samba/idmap/hash.so
/usr/lib/samba/idmap/ldap.so
/usr/lib/samba/idmap/tdb2.so
/usr/sbin
/usr/sbin/winbindd
/usr/bin
/usr/bin/wbinfo
/usr/bin/ntlm_auth
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/wbinfo.1.gz
/usr/share/man/man1/ntlm_auth.1.gz
/usr/share/man/man8
/usr/share/man/man8/winbindd.8.gz
/usr/share/man/man8/idmap_ldap.8.gz
/usr/share/man/man8/idmap_nss.8.gz
/usr/share/man/man8/idmap_rid.8.gz
/usr/share/man/man8/idmap_tdb.8.gz
/usr/share/man/man8/idmap_tdb2.8.gz
/usr/share/man/man8/idmap_hash.8.gz
/usr/share/man/man8/pam_winbind.8.gz
/usr/share/man/man8/idmap_ad.8.gz
/usr/share/man/man8/idmap_adex.8.gz
/usr/share/pam-configs
/usr/share/pam-configs/winbind
/usr/share/doc
/usr/share/doc/winbind
/usr/share/doc/winbind/copyright
/usr/share/doc/winbind/NEWS.Debian.gz
/usr/share/doc/winbind/README.build.gz
/usr/share/doc/winbind/changelog.Debian.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/winbind
/lib
/lib/security
/lib/security/pam_winbind.so
/lib/libnss_winbind.so.2
/lib/libnss_wins.so.2
/etc
/etc/logrotate.d
/etc/logrotate.d/winbind
/etc/default
/etc/default/winbind
/etc/init.d
/etc/init.d/winbind
/usr/lib/samba/nss_info/sfu20.so
/usr/lib/samba/nss_info/sfu.so
/usr/lib/samba/nss_info/rfc2307.so
